# root/lib/systemd/system/clamscan-nethgui.service
[Unit]
Description=Start a scan in nethgui
Requisite=httpd-admin.service

[Service]
Type=simple
ExecStart=/sbin/e-smith/nethserver-clamscan
